# test_scenario/unitfiles/1-2.target
[Install]
WantedBy = default.target

[Unit]
Description= "Target that can be referenced by services so they get pulled in"
Before=2-3.target