# test_scenario/unitfiles/1_1.service
[Service]
ExecStart= /usr/bin/echo "1_1"

[Install]
WantedBy=1-2.target