# test_scenario/unitfiles/1_2.service
[Unit]
Before=1-2.target

[Service]
ExecStart= /bin/sleep 1
Type=oneshot

[Install]
WantedBy=1-2.target