# test_scenario/unitfiles/2_1.service
[Unit]
After=1-2.target

[Service]
ExecStart= /usr/bin/echo "2_1"

[Install]
WantedBy=default.target