# test_scenario/unitfiles/2_2.service
[Unit]
After=1-2.target

[Service]
ExecStart= /usr/bin/echo "2_2"
Type=oneshot

[Install]
WantedBy=2-3.target