# test_scenario/unitfiles/3_1.service
[Unit]
After=2-3.target

[Service]
ExecStart= /usr/bin/echo "3_1"
Type=oneshot

[Install]
WantedBy=default.target