# test_scenario/unitfiles/3_2.service
[Unit]
After=2-3.target

[Service]
ExecStart= /usr/bin/echo "3_2"
Type=oneshot

[Install]
WantedBy=default.target