# test_scenario/unitfiles/3_3.service
[Unit]
After=2-3.target

[Service]
ExecStart= /bin/sleep 1
Type=oneshot

[Install]
WantedBy=default.target