# test_scenario/unitfiles/not_loaded.service
[Unit]

[Service]
ExecStart= /usr/bin/echo "Do not execute this shouldnt even be loaded"
Type=oneshot

[Install]
