# test_scenario/unitfiles/sockets.target
[Unit]
Description= Target that synchronizes startup. At this point all sockets should be up

[Install]
WantedBy=default.target