# test_scenario/unitfiles/unloaded.socket
[Unit]
Description=test Socket

[Socket]
FileDescriptorName= TestSockUnloaded

[Install]
WantedBy=sockets.target