# test_units_disabled/ls.service
[Unit]
requires= ls2
Before= ls2

[Service]
Exec= /usr/bin/ls /etc
Restart= no