# test_units_disabled/ls2.service
[Service]
Exec= /usr/bin/ls /etc
Restart= no

[Install]
WantedBy= ls