# test_units_disabled/ls2.service
[Service]
Exec= /usr/bin/ls /etc
keep_alive = false

[Install]
WantedBy= ls