# test_units_disabled/ls3.service
[Unit]
requires= ls4
Before= ls4

[Service]
Exec= /usr/bin/sleep 100
Stop= killall sleep
Restart= no