# test_units_disabled/ls4.service
[Service]
Exec= /usr/bin/sleep 2
Restart= no

[Install]
WantedBy= ls3