# test_units_disabled/service_using_dbus.service
[Service]
Exec= /bin/true
Type= dbus
BusName= org.freedesktop.DBus
Restart= no