# test_units/ls.service
[Unit]
requires= ls2
Before= ls2

[Service]
Exec= ls /etc
keep_alive = false