# test_units/ls.service
[Service]
Exec= ls /etc
keep_alive = true