# test_units/ls2.service
[Service]
Exec= ls /usr
keep_alive = true