# test_units/ls3.service
[Unit]
requires= ls4
Before= ls4

[Service]
Exec= sleep 2
keep_alive = true