# test_units/ls4.service
[Service]
Exec= sleep 2
keep_alive = false

[Install]
WantedBy= ls3