# test_units/ls4.service
[Service]
Exec= sleep 2
keep_alive = true

[Install]
WantedBy= ls3