# test_units/oneshot.service
[Service]
ExecStart= /bin/sleep 3 
ExecStartPre= /bin/sleep 1
ExecStartPost= /bin/sleep 1
ExecStop = /bin/sleep 1
Type= oneshot
TimeoutStartSec= 0min 1s

[Unit]
Before= test2.service

[Install]
WantedBy= default.target