# test_units/oneshot.service
[Service]
ExecStart= /bin/false
ExecStartPre= /bin/sleep 1
ExecStartPost= /bin/sleep 1
ExecStop = echo "Stop"
ExecStop = /bin/false
ExecStopPost = echo "Post stop"
Type= oneshot
TimeoutStartSec= 0min 4s

[Unit]
Before= test2.service

[Install]
WantedBy= default.target