# test_units/service_using_dbus.service
[Service]
Exec= /bin/true
Type= dbus
BusName= "org.freedesktop.DBus"
keep_alive = false