# test_units/test2.service
[Service]
Exec= ./target/debug/testservice arg1 arg2
Sockets= test2
NotifyAccess= all
Type= notify
keep_alive = true

[Install]
Before= test
