# test_units/test.service
[Service]
Exec= ./target/debug/testservice
keep_alive = false
Sockets= test_second
