# test_units/test.service
[Service]
ExecStart= ./target/debug/testservice arg1 arg2
Sockets= test_second.socket
NotifyAccess= all
Type= notify
Restart= always
ExecStartPre= /bin/sleep 1
ExecStartPre= /bin/echo "This should be printed with the correct prefix even tough it is a message from ExecStartPre"
ExecStartPost= /bin/sleep 1
ExecStop = /bin/sleep 1
ExecStop = /bin/sleep 1

[Unit]
Before= test2.service

[Install]
WantedBy= default.target