# test_units/test.service
[Service]
ExecStart= ./target/debug/testservice arg1 arg2 "this is a multi word argument"
Sockets= test_second.socket
NotifyAccess= all
Type= notify
Restart= always
ExecStartPre= echo "perstart"
ExecStartPre= echo "This should be printed with the correct prefix even tough it is a message from ExecStartPre"
ExecStartPost= echo "posstart"
ExecStop = echo "stop"
ExecStopPost = echo "poststop"

[Unit]
Before= test2.service

[Install]
WantedBy= default.target