# test_units/test.service
[Service]
Exec= ./test_services/test.sh
keep_alive = false