# test_units/test.socket
[Unit]
Description=test Socket
PartOf=test.service

[Socket]
File= ./logs/servicelog

[Install]
WantedBy=sockets.target