# test_units/test.service
[Service]
Exec= ./target/debug/testservice arg1 arg2
Sockets= test_second.socket
NotifyAccess= all
Type= notify
keep_alive = true

[Unit]
Before= test2.service

[Install]
WantedBy= default.target