# test_units/test.socket
[Unit]
Description=test Socket
PartOf=test.service

[Socket]
ListenStream= ./servicelog

[Install]
WantedBy=sockets.target