# test_units/test.service
[Service]
Exec= ./target/debug/testservice
Sockets= test_second
NotifyAccess= all
Type= notify
keep_alive = true