# test_units/unloaded.service
[Service]
ExecStart= ./target/debug/testservice arg1 arg2 "this is a multi word argument"

[Unit]

[Install]
