# test_units/unloaded.socket
[Unit]
Description=test Socket

[Socket]
FileDescriptorName= TestSockUnloaded

# ListenStream= 127.0.0.1:8090
[Install]
WantedBy=sockets.target