# example_boot_units/default.target
[Unit]
Description= Currently starts you at a agetty instance at tty1