# example_boot_units/getty.target
[Unit]
Description= All getty instances have been started before this target is reached

[Install]
RequiredBy=default.target