# example_boot_units/mount_all.service
[Unit]
Description= Mount all files in /etc/fstab

[Service]
Type= oneshot
ExecStart= /sbin/mount -a


[Install]
WantedBy= local_filesystem.target